# etc/crontab
CONSOLE={absolute path to ./bin/console}

* */2 * * * $CONSOLE document:index:update